# speed_test.service
# path:   /home/klassiker/.local/share/repos/python/speed_test.service
# author: klassiker [mrdotx]
# github: https://github.com/mrdotx/python
# date:   2021-01-15T13:51:41+0100

[Unit]
Description=speed test in user contaxt alarm

[Service]
Type=oneshot
ExecStart=%h/.local/share/repos/python/speed_test.py

[Install]
WantedBy=multi-user.target
